# scripts/nanoclaw.service
[Unit]
Description=NanoClaw
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/home/exedev/nanoclaw
ExecStart=/home/exedev/.local/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=PATH=/home/exedev/.local/bin:/usr/local/bin:/usr/bin:/bin
User=exedev

[Install]
WantedBy=multi-user.target
